# src/launchers.crontab
*/5  * * * *  bash /home/selknam/dev/rata/src/rates_launcher.bash 5 >> /home/selknam/var/log/rates_launcher.5.log  2>&1
06   * * * *  nice -n 0 bash /home/selknam/dev/rata/src/launch_prepare_and_model_binary_classifiers.bash >> /home/selknam/var/log/launch_prepare_and_model_binary_classifiers.log 2>&1
*/5  * * * *  nice -n 0 bash /home/selknam/dev/rata/src/launch_forecast_binary_classifiers.bash >> /home/selknam/var/log/launch_forecast_binary_classifiers.log 2>&1
#